# config/selfstarter.conf
[program:selfstarter]
command=rails s
directory=/home/vagrant/selfstarter
user=vagrant
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
